# files/system_files/deck/usr/lib/systemd/system/zelideck-autologin.service
[Unit]
Description=Enables ZeliDeck autologin
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/zelideck-autologin

[Install]
WantedBy=multi-user.target
