# config/files/deck/lib/systemd/system/return-to-gamemode.service
[Unit]
Description=Return to Game Mode

[Service]
Type=oneshot
ExecStart=/usr/bin/return-to-gamemode