# system_files/deck/usr/lib/systemd/user/zelideck-user-setup.service
[Unit]
Description=Configure Zeliblue Deck for current user

[Service]
Type=simple
ExecStart=/usr/libexec/zelideck-user-setup

[Install]
WantedBy=default.target
