# files/system_files/shared/etc/profile.d/zeliblue-fastfetch.sh
#!/bin/sh

alias fastfetch='fastfetch --logo fedora-silverblue'
