# system_files/shared/etc/profile.d/zeliblue-fastfetch.sh
#!/bin/sh

alias fastfetch="/usr/bin/fastfetch --logo fedora-silverblue -c /usr/share/zeliblue/fastfetch/zeliblue.jsonc"
