# config/files/deck/usr/lib/systemd/system/zelideck-toggle-autologin.service
[Unit]
Description=Toggle desktop autologin

[Service]
Type=oneshot
ExecStart=/usr/bin/zelideck-toggle-autologin